# root/etc/services.d/mongodb/run
#!/usr/bin/with-contenv bash

exec \
        s6-setuidgid abc /usr/bin/mongod -f /etc/mongodb.conf --fork --logpath /var/log/mongodb.log --config /config/mongodb.conf